# deploy/study-seat.service
[Unit]
Description=Study Seat Reservation System
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/study-seat/backend
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=ANTHROPIC_AUTH_TOKEN=sk-6a52c8307b45464c84ba19e1c0d4a1ad
Environment=ANTHROPIC_BASE_URL=https://api.deepseek.com/anthropic
Environment=ANTHROPIC_MODEL=deepseek-v4-pro

[Install]
WantedBy=multi-user.target
